# resources/nodered.socket
[Unit]
Description=Node-RED Socket
PartOf=nodered.service

[Socket]
ListenStream=1880

[Install]
WantedBy=sockets.target
